# supervisor-investorservitude.conf
[program:investorservitude]
command=/home/tim/gunicorn.sh
directory=/home/tim
user=tim
autostart=true
autorestart=true
stdout_logfile = /var/log/supervisor/investorservitude.log
redirect_stderr=true